# root/etc/systemd/system/NetworkManager.service.d/NetworkManager-load-state.conf
[Service]
ExecStart=
ExecStart=/usr/sbin/NetworkManager --no-daemon --state-file=/persist/NetworkManager/NetworkManager.state
